# bin/unofficialtranslations.service
[Unit]
Description="unofficialtranslations HTTP server"

[Service]
ExecStart=/usr/bin/node /home/ubuntu/unofficialtranslations.com/bin/www.js
#WorkingDirectory=/home/ubuntu/unofficialtranslations.com/
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifies=UnofficialTranslations
Environment=NODE_ENV=production PORT=3000

[Install]
WantedBy=multi-user.target